# com.scalapay.deal-qualifier.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.scalapay.deal-qualifier</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/stefano.conforti@scalapay.com/Cursor/sales-qualifier/run_agent.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/stefano.conforti@scalapay.com/Cursor/sales-qualifier/agent.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/stefano.conforti@scalapay.com/Cursor/sales-qualifier/agent.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HUBSPOT_TOKEN</key>
        <string>YOUR_HUBSPOT_TOKEN</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
